# conf/systemd.timer
[Unit]
Description=Run __APP__ regularly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target